# systemd/.config/systemd/user/conky.service
# -*- mode: conf; -*-

[Unit]
Description=conky

[Service]
Type=forking
ExecStart=/home/ht/bin/start-conky
Environment="DISPLAY=:0"
