# xdg/.config/autostart/xrdb.desktop
[Desktop Entry]
Name=xrdb
Type=Application
Exec=/usr/bin/xrdb -load /home/ht/.Xresources
Terminal=false
